# deploy/meddcode.service
[Unit]
Description=MedDcode investor demo app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/meddcode
Environment=PORT=8501
Environment=HOST=0.0.0.0
Environment=MEDDCODE_DIST_DIR=/home/opc/meddcode/dist
Environment=OCI_GENAI_API_KEY_PATH=/home/opc/.keys/api_genai_key1
Environment=OCI_REGION=us-ashburn-1
ExecStart=/usr/bin/node /home/opc/meddcode/server.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
